# src/daemon/plasmazones.service.in
[Unit]
Description=PlasmaZones Window Snapping Daemon
Documentation=https://github.com/fuddlesworth/PlasmaZones
PartOf=graphical-session.target
After=plasma-workspace.target plasma-kglobalaccel.service

[Service]
Type=simple
ExecStart=@KDE_INSTALL_FULL_BINDIR@/plasmazonesd
Restart=on-failure
RestartSec=5
Slice=session.slice
Environment=QT_PLUGIN_PATH=@KDE_INSTALL_FULL_LIBDIR@/qt6/plugins
Environment=LD_LIBRARY_PATH=@KDE_INSTALL_FULL_LIBDIR@

[Install]
WantedBy=plasma-workspace.target
