# packaging/arch/plasmazones.install
# PlasmaZones pacman install script
# SPDX-License-Identifier: GPL-3.0-or-later

post_install() {
    echo ""
    echo "══════════════════════════════════════════════════════════════════"
    echo "  PlasmaZones installed successfully!"
    echo "══════════════════════════════════════════════════════════════════"
    echo ""
    echo "  To enable the daemon:"
    echo "      systemctl --user enable --now plasmazones.service"
    echo ""
    echo "  The KWin effect is enabled by default, but KWin must be"
    echo "  restarted to load it. Log out and back in, or run:"
    echo "      kwin_wayland --replace &"
    echo ""
    echo "  Settings: System Settings → Window Management → PlasmaZones"
    echo ""
    echo "  If PlasmaZones doesn't appear in System Settings, run:"
    echo "      kbuildsycoca6 --noincremental"
    echo "  or log out and back in."
    echo ""
    echo "  Snap Assist thumbnails require a fresh session to pick up the"
    echo "  screenshot auth. Log out and back in after install."
    echo ""
}

post_upgrade() {
    echo ""
    echo "══════════════════════════════════════════════════════════════════"
    echo "  PlasmaZones upgraded!"
    echo "══════════════════════════════════════════════════════════════════"
    echo ""
    echo "  If the daemon was running, restart it:"
    echo "      systemctl --user restart plasmazones.service"
    echo ""
    echo "  Note: KWin must be restarted to pick up effect plugin changes."
    echo "  Log out and back in, or run: kwin_wayland --replace &"
    echo ""
}

post_remove() {
    echo ""
    echo "  To fully remove PlasmaZones:"
    echo "      systemctl --user disable --now plasmazones.service"
    echo ""
}
